# wgapi.service
[Unit]
Description=Wireguard API
After=network.target

[Service]
ExecStart=/usr/local/bin/wgapi

[Install]
WantedBy=multi-user.target
